# etc/cron.weekly/vkpurge-rm-all.sh
#!/bin/sh

if command -v vkpurge >/dev/null; then
  vkpurge rm all
fi
