# etc/cron.weekly/xbps-remove-clean-cache.sh
#!/bin/sh

if command -v xbps-remove >/dev/null; then
  xbps-remove -O
fi
