# system/root.crontab
* 00-15 * * * chmod -x /usr/bin/prismlauncher
* 16-23 * * * chmod +x /usr/bin/prismlauncher
