# system/cron.weekly/remove
#!/bin/sh
xbps-remove -O
vkpurge rm all
