# system/cron.weekly/remove
#!/bin/sh
xbps-remove -Oyo
vkpurge rm all
