# system/cron.weekly/vacuum
#!/bin/sh
cd /var/log/socklog && rm */@*
cd ~/.local/share && rm -fr Trash/*
