# system/cron.monthly/retain
#!/bin/sh
borg prune -m 12 -w 4 -d 7 /borg
borg compact /borg
