# system/cron.monthly/vacuum
#!/bin/sh
xbps-remove -O
fstrim /
rm ~/.local/share/Trash/*/*
rm /var/log/socklog/*/@*
