# system/cron.daily/backup
#!/bin/sh
borg create "$@" -xC 'auto,zstd,18' /borg::"$(date '+%Y%m%d')" / /boot
