# system/cron.daily/update
#!/bin/sh
xbps-install -Syu
vkpurge rm all
