# system/cron.daily/update
#!/bin/sh
xbps-install -Syu xbps
xbps-install -Syu
